# crates/chronos-sandbox/quadlets/python-target.container
[Unit]
Description=Chronos sandbox Python target

[Container]
Image=python:3.12
ContainerName=chronos-sandbox-python
AddCapability=SYS_PTRACE
SecurityLabelDisable=true
Annotation=run.oci.keep_original_groups=1

[Service]
Restart=no
