# crates/chronos-sandbox/quadlets/rust-target.container
[Unit]
Description=Chronos sandbox Rust target

[Container]
Image=rust:1.77
ContainerName=chronos-sandbox-rust
AddCapability=SYS_PTRACE
SecurityLabelDisable=true
Annotation=run.oci.keep_original_groups=1

[Service]
Restart=no
